# timed_network_save_push.timer
[Unit]
Description=Timer for network save push

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
